# arch/.config/systemd/user/bkup.service
[Unit]
Description=Updates backup for /home/USER

[Service]
Type=oneshot
ExecStart=/home/USER/bin/bkup

[Install]
WantedBy=default.target
