# arch/.config/systemd/user/bkup.timer
[Unit]
Description=Send backup of /home/USER

[Timer]
OnBootSec=15m
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
