# arch/.config/systemd/user/feh.timer
[Unit]
Description=Change background every 10 minutes

[Timer]
OnActiveSec=5sec
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
