# arch/.config/systemd/user/feh.service
[Unit]
Description=Changes background

[Service]
Type=oneshot
ExecStart=/home/dimus/bin/nitro

[Install]
WantedBy=default.target
